# deploy/systemd/roboktober-deploy.service
[Unit]
Description=Roboktober deployment job
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Set to the user that owns /var/www/www.roboktober.nl
User=rein
Group=www-data
WorkingDirectory=/var/www/www.roboktober.nl
Environment=REPO_DIR=/var/www/www.roboktober.nl
Environment=BRANCH=master
Environment=RUN_MIGRATIONS=true
Environment=BUILD_FRONTEND=false
Environment="PHP_CANDIDATES=php8.6 php8.5 php8.4 php8.3 php"
ExecStart=/usr/bin/env bash /var/www/www.roboktober.nl/deploy/deploy.sh

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
